# system_files/etc/systemd/user/llm-server-update.service
[Unit]
Description=Update LLM server

[Service]
Type=oneshot
ExecStart=sh -c 'llm-server.sh update'
