# system_files/etc/systemd/user/llm-server-update.timer
[Unit]
Description=Daily LLM server update

[Timer]
OnBootSec=5min-10min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
