# system_files_headless/etc/systemd/user/dyndns.timer
[Unit]
Description=Make us known to the public

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
