# system_files_headless/etc/systemd/system/ausgabenzettel-update.timer
[Unit]
Description=Checks for ausgabenzettel images

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
