# system_files_headless/etc/systemd/system/postgres-update.timer
[Unit]
Description=Check for new Postgres image daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
